# local/etc/cron.d/logrotate
* * * * * root /etc/cron.daily/logrotate
