# packagers/deb/reactor-manager/DEBIAN/postinst
#!/bin/bash

set -e

function service_start() {
    invoke-rc.d reactor-manager start || exit $?
}

case $1 in
    configure) service_start;;
esac
