# packagers/deb/reactor-manager/DEBIAN/prerm
#!/bin/bash

set -e

invoke-rc.d reactor-manager stop || exit $?
