# packagers/deb/reactor-server/DEBIAN/postinst
#!/bin/bash

set -e

function service_start() {
    invoke-rc.d reactor start || exit $?
}

case $1 in
    configure) service_start;;
esac
