# packagers/deb/reactor-server/DEBIAN/prerm
#!/bin/bash

set -e

invoke-rc.d reactor-server stop || exit $?
